# automagik-forge.service
[Unit]
Description=Automagik Forge - Task Management & AI Orchestration
After=network.target
Documentation=https://github.com/namastexlabs/automagik-forge

[Service]
Type=forking
User=produser
WorkingDirectory=/home/produser/prod/automagik-forge
Environment="PATH=/home/produser/.nvm/versions/node/v22.21.1/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PM2_HOME=/home/produser/.pm2"

# Start PM2 with local ecosystem config
ExecStart=/home/produser/.nvm/versions/node/v22.21.1/bin/pm2 start ecosystem.config.cjs
ExecStop=/home/produser/.nvm/versions/node/v22.21.1/bin/pm2 stop ecosystem.config.cjs
ExecReload=/home/produser/.nvm/versions/node/v22.21.1/bin/pm2 restart ecosystem.config.cjs

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=400
StartLimitBurst=3

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
